# deploy/safeplay-agent.service
[Unit]
Description=SafePlay AI Agent - Self-healing YouTube Downloader Agent
Documentation=https://github.com/matbram/safeplay-yt-dlp
After=network.target safeplay-ytdlp.service
Wants=safeplay-ytdlp.service

[Service]
Type=simple
User=safeplay
Group=safeplay
WorkingDirectory=/opt/safeplay-ytdlp

# Environment
EnvironmentFile=/opt/safeplay-ytdlp/.env
Environment=PYTHONUNBUFFERED=1

# Start the agent
ExecStart=/opt/safeplay-ytdlp/venv/bin/python -m agent.main

# Restart configuration
Restart=always
RestartSec=10

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=safeplay-agent

# Security hardening
NoNewPrivileges=false
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
